# lempirahol.io.kittntabbr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>lempirahol.io.kittntabbr</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lempiraholio/Development/KittnTabbR/KittnTabbR.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lempiraholio/Development/KittnTabbR</string>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/lempiraholio/Library/Logs/kittntabbr.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lempiraholio/Library/Logs/kittntabbr.log</string>
</dict>
</plist>
